# docs/gentoo-init
#!/sbin/runscript

# for Gentoo; store as /etc/init.d/amqpprocessmanager
# configure variables below
# also configure PIDFILE, CONFIG and LOG_CONFIG in /etc/conf.d/amqpprocessmanager

jar=/path/to/assembled.jar

depend() {
        after net
        use dns
}

start() {
        ebegin "Starting AMQP Process Manager"
        start-stop-daemon --start --quiet --background --make-pidfile --pidfile "${PIDFILE}" -u REDACTED --exec "$(java-config -J)" -- -Djava.util.logging.config.file="${LOG_CONFIG}" -jar "${jar}" "${CONFIG}"
        eend $?
}

stop() {
        ebegin "Stopping AMQP Process Manager"
        start-stop-daemon --stop --quiet --pidfile "${PIDFILE}"
        eend $?
}

restart() {
        stop
        start
}
